# skills/meta-skills/systemd/meta-skills-daily.timer.in
[Unit]
Description=Run meta-skills daily at @DESCRIPTION_TIME@

[Timer]
OnCalendar=@ON_CALENDAR@
Persistent=true

[Install]
WantedBy=timers.target
